# supervisord.conf
[program:vncserver]
command=vncserver :1 -geometry 1600x900 -depth 24
environment=USER=bketelsen

[supervisord]
logfile=/home/bketelsen/supervisord.log

